# debian/xinetd.d/amanda
# default: on
# description: The amanda service
service amanda
{
#	only_from	= <Amanda server>
	socket_type	= dgram
	protocol	= udp
	wait		= yes
	user		= backup
	group		= backup
	groups		= yes
	server		= /usr/lib/amanda/amandad
	server_args	= -auth=bsd amdump amindexd amidxtaped
	disable		= no
}
