# debian/xinetd.d/amandaidx
#default: on
# description: The amanda index service
service amandaidx
{
#	only_from	= <Amanda client list>
	socket_type	= stream
	protocol	= tcp
	wait		= no
	user		= backup
	group		= backup
	groups		= yes
	server		= /usr/lib/amanda/amindexd
	server_args	= -auth=bsd amdump amindexd amidxtaped
	disable		= no
}
